# scripts/com.cactus.docker-cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cactus.docker-cleanup</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/prueba/Desktop/CactusDashboard/scripts/docker-cleanup.sh</string>
        <string>--force</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/prueba/Desktop/CactusDashboard/logs/docker-cleanup.out</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/prueba/Desktop/CactusDashboard/logs/docker-cleanup.err</string>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>WorkingDirectory</key>
    <string>/Users/prueba/Desktop/CactusDashboard</string>
</dict>
</plist>
